# dot_config/systemd/user/flameshot.service
[Unit]
Description=flameshot

[Service]
Type=simple
Restart=always
ExecStart=flameshot

[Install]
WantedBy=default.target
