# dot_config/systemd/user/update-wallpaper.service
[Unit]
Description=Update wallpaper using custom script

[Service]
Type=oneshot
ExecStart=%h/bin/updatewal
Slice=background-graphical.slice
