# dot_config/systemd/user/xremap.service
[Unit]
Description=xremap

[Service]
Type=simple
Restart=always
ExecStart=xremap /home/bingo/.config/xremap/config.yml --watch

[Install]
WantedBy=default.target
