# dot_config/systemd/user/blueman-applet.service.d/override.conf
[Unit]
After=graphical-session.target
Requisite=graphical-session.target

[Service]
Slice=app-graphical.slice
Restart=on-failure

[Install]
WantedBy=graphical-session.target
